# packaging/systemd/sage-gloved.service.in
[Unit]
Description=Sage Glove agent-session supervisor
Documentation=https://github.com/sage-protocol/sage-protocol

[Service]
Type=simple
ExecStart="@GLOVED_BINARY@" --config "@GLOVE_CONFIG@"
Restart=on-failure
RestartSec=2s
TimeoutStopSec=15s
UMask=0077

[Install]
WantedBy=default.target
