# packaging/launchd/org.sage-protocol.gloved.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.sage-protocol.gloved</string>
    <key>ProgramArguments</key>
    <array>
        <string>@GLOVED_BINARY@</string>
        <string>--config</string>
        <string>@GLOVE_CONFIG@</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>ThrottleInterval</key>
    <integer>2</integer>
</dict>
</plist>
